# eos-setup-registry/eos-setup-registry/usr/lib/systemd/system/registry-rw.service
[Unit]
Description=Activate Container Read/Write Registry
Before=registry_resize.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/sbin/set_rw_registry

[Install]
WantedBy=multi-user.target
