# eos-setup-os-update/eos-setup-os-update/usr/lib/systemd/system/os-fetch.service
[Unit]
Description=OS Fetch Images

[Service]
Type=oneshot
ExecStart=/srv/www/fleet/os-images/fetch

[Install]
WantedBy=multi-user.target
