# eos-setup-registry/eos-setup-registry/eos-setup-registry-1.1.1/eos-setup-registry/usr/lib/systemd/system/registry_resize.service
[Unit]
Description=Resize Container Registry
After=registry-rw.service
Requires=registry-rw.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/registry_resize

[Install]
WantedBy=multi-user.target
