# eos-setup-os-update/eos-setup-os-update/eos-setup-os-update-1.1.1/eos-setup-os-update/usr/lib/systemd/system/os-update-daemon@.timer
[Unit]
Description=OS update daemon periodically for %i.raw
Requires=network-online.target
After=network.target syslog.target

[Timer]
Unit=os-update-daemon@%i.service
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=multi-user.target
