# eos-setup-os-update/eos-setup-os-update/eos-setup-os-update-1.1.1/eos-setup-os-update/usr/lib/systemd/system/os-update-daemon@.service
[Unit]
Description=OS update daemon for %i.raw

[Service]
Type=simple
EnvironmentFile=/etc/os-update-daemon.conf
ExecStart=/usr/bin/os-update-daemon.sh --image /srv/www/fleet/os-images/%i.raw --oneshot
