# eos-setup-os-update/eos-setup-os-update/eos-setup-os-update-1.1.1/eos-setup-os-update/usr/lib/systemd/system/update_commit.service
[Unit]
Description=Commit OS Update

[Service]
Type=oneshot
ExecStart=/usr/sbin/os-update --commit --reboot

[Install]
WantedBy=multi-user.target
